# supervisord-solr.conf
[program:solr]
command=/start-solr.sh
numprocs=1
autostart=true
autorestart=true
