# integrationtestsettings/supervisord-fuseki.conf
[program:fuseki]
command=/start-fuseki.sh
numprocs=1
autostart=true
autorestart=true
